# rc.local
#!/bin/sh -e
#
# rc.local
#
# This script is executed at the end of each multiuser runlevel.
# Make sure that the script will "exit 0" on success


/home/pi/ip/ip.sh &	# Lanza script de envio de ip por mail
/home/pi/ip/ipmail.sh & # Lanza el script de IP hour check


exit 0
